# spark-cluster/docker-compose/pymaster/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:spark_master]
command=/usr/local/spark/sbin/start-master.sh